# buildroot-external/rootfs-overlay/lib/udev/rules.d/02-dp_hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/bin/pa-udev-update-sink.sh"
